# docker/al2023-systemd/systemd-email@.service
[Unit]
Description=Send status email for %i

[Service]
Type=oneshot
ExecStart=/etc/systemd-email.sh systemd %i
User=nobody
Group=systemd-journal